# templates/supervisord/program.conf.erb
[program:<%= @name %>]
command=<%= @command %>
user=<%= @user %>
autostart=true
autorestart=true
stdout_logfile=<%= @stdout_file %>
stderr_logfile=<%= @stderr_file %>